# .contrib/bettertunnels-caddy.service
[Unit]
Description=BetterTunnels — Caddy reverse proxy (tunnels.local:80 → localhost:5730)
After=bettertunnels.service
BindsTo=bettertunnels.service

[Service]
Type=exec
ExecStart=/usr/bin/caddy run --config /home/shahriyar/dev/next/cloudtunnel/.contrib/Caddyfile
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
